# blews.service
[Unit]
Description=Blews Service

[Service]
Type=simple
ExecStart=/var/blews/blews-backend/start-server.sh

[Install]
WantedBy=multi-user.target